# ops/runtime/systemd/bridge-watchdog.service
[Unit]
Description=Pictronic Runtime/Bridge Watchdog
After=network.target pictronic-runtime.service bridge-connector.service

[Service]
Type=simple
WorkingDirectory=/home/mataleao/.paperclip/instances/default/projects/fb3b0db7-d166-45d8-93b9-552802652d7d/4bbbf7cb-906f-4da9-93f5-189a4de11489/_default
Environment=PICTRONIC_RUNTIME_ENV_FILE=/home/mataleao/.paperclip/instances/default/projects/fb3b0db7-d166-45d8-93b9-552802652d7d/4bbbf7cb-906f-4da9-93f5-189a4de11489/_default/.env.runtime
Environment=PICTRONIC_PROJECT_ROOT=/home/mataleao/.paperclip/instances/default/projects/fb3b0db7-d166-45d8-93b9-552802652d7d/4bbbf7cb-906f-4da9-93f5-189a4de11489/_default
EnvironmentFile=-/home/mataleao/.paperclip/instances/default/projects/fb3b0db7-d166-45d8-93b9-552802652d7d/4bbbf7cb-906f-4da9-93f5-189a4de11489/_default/.env.runtime
ExecStart=/usr/bin/python3 /home/mataleao/.paperclip/instances/default/projects/fb3b0db7-d166-45d8-93b9-552802652d7d/4bbbf7cb-906f-4da9-93f5-189a4de11489/_default/ops/runtime/watchdog.py --manager systemd --runtime-service pictronic-runtime --connector-service bridge-connector --env-file ${PICTRONIC_RUNTIME_ENV_FILE} --runtime-error-log docs/e2e/pictronic-runtime.systemd.err.log --connector-error-log docs/e2e/pictronic-connector.systemd.err.log --poll-interval-seconds 10 --recovery-log docs/e2e/pictronic-watchdog.systemd-recovery.log --state-file docs/e2e/pictronic-watchdog.systemd.state.json
Restart=always
RestartSec=2
StandardOutput=append:/home/mataleao/.paperclip/instances/default/projects/fb3b0db7-d166-45d8-93b9-552802652d7d/4bbbf7cb-906f-4da9-93f5-189a4de11489/_default/docs/e2e/pictronic-watchdog.systemd.log
StandardError=append:/home/mataleao/.paperclip/instances/default/projects/fb3b0db7-d166-45d8-93b9-552802652d7d/4bbbf7cb-906f-4da9-93f5-189a4de11489/_default/docs/e2e/pictronic-watchdog.systemd.err.log

[Install]
WantedBy=multi-user.target
